# stock-system/scripts/systemd/stock_worker.service
[Unit]
Description=Stock System Scheduler Worker
After=network.target

[Service]
EnvironmentFile=/etc/stock-system.env
WorkingDirectory=/opt/stock-system
ExecStart=/opt/stock-system/.venv/bin/python -m app.scheduler.cli --once daily
Restart=on-failure
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
