# stock-system/scripts/systemd/stock_api.service
[Unit]
Description=Stock System API
After=network.target

[Service]
EnvironmentFile=/etc/stock-system.env
WorkingDirectory=/opt/stock-system
ExecStart=/opt/stock-system/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000 --workers 1 --loop uvloop --http h11
Restart=on-failure
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
